# mcproxy.service
[Unit]
Description=MCProxy MeshCom Message Proxy
After=network.target mcproxy-ble.service
Wants=mcproxy-ble.service

[Service]
Type=simple
User=martin
WorkingDirectory=/home/martin/mcproxy-test
Environment=PATH=/home/martin/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=MCADVCHAT_ENV=dev
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/martin/.local/bin/uv run mcproxy
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
