# bootstrap/templates/mcapp.service
[Unit]
Description=McApp - MeshCom Message Proxy
Documentation=https://github.com/DK5EN/McAdvChat
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{HOME}}/mcapp
Environment="PATH={{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="MCAPP_ENV=prod"
ExecStart={{HOME}}/.local/bin/uv run mcapp
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/log/mcapp /etc/mcapp /var/lib/mcapp {{HOME}}/mcapp {{HOME}}/.cache/uv

[Install]
WantedBy=multi-user.target
