# bootstrap/templates/mcapp-update.path
[Unit]
Description=Watch for McApp update trigger

[Path]
PathExists=/var/lib/mcapp/update-trigger

[Install]
WantedBy=multi-user.target
