# bootstrap/templates/mcproxy.service
[Unit]
Description=MCProxy - MeshCom Message Proxy
Documentation=https://github.com/DK5EN/McAdvChat
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory=/usr/local/bin
Environment="PATH={{HOME}}/mcproxy-venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="MCADVCHAT_ENV=prod"
ExecStart={{HOME}}/mcproxy-venv/bin/python /usr/local/bin/C2-mc-ws.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/log/mcproxy /etc/mcadvchat

[Install]
WantedBy=multi-user.target
